# systemd/dashboard-refresh.timer
[Unit]
Description=Run DBA dashboard refresh every 15 minutes (Plan 38)

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
Unit=dashboard-refresh.service

[Install]
WantedBy=timers.target
