# systemd/source-health-permits.timer
[Unit]
Description=Somerville permits hourly source-health check schedule
Requires=source-health-permits.service

[Timer]
OnCalendar=*-*-* *:10:00
Persistent=true

[Install]
WantedBy=timers.target
